# monitoring/service-checks/crontab
# Hourly checks at offset minutes to avoid simultaneity
17 * * * *  python3 /checks/iptv-auth.py            >> /proc/1/fd/1 2>&1
23 * * * *  python3 /checks/iptv-channel-count.py   >> /proc/1/fd/1 2>&1
29 * * * *  python3 /checks/iptv-epg-freshness.py   >> /proc/1/fd/1 2>&1
35 * * * *  python3 /checks/iptv-canary-epg.py     >> /proc/1/fd/1 2>&1
# Daily renewal warning at 09:00 local
0 9 * * *   python3 /checks/iptv-renewal-warn.py    >> /proc/1/fd/1 2>&1

# Every 5 minutes: probe each public dephekt.net URL from outside the docker
# network. Catches proxy/tunnel/edge breakage that container healthchecks and
# newt's per-target healthchecks miss (both probe inside the docker network).
*/5 * * * * python3 /checks/public-availability.py >> /proc/1/fd/1 2>&1
